# scripts/si-production-bot.service
[Unit]
Description=SI Production Telegram Bot
After=network.target

[Service]
Type=simple
User=ae
WorkingDirectory=/home/ae/si-production
Environment="NODE_ENV=production"
Environment="NVM_DIR=/home/ae/.nvm"
ExecStart=/bin/bash -c "source /home/ae/.nvm/nvm.sh && cd /home/ae/si-production && node scripts/start-bot-polling.js"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
